# packages/coach/payload/etc/systemd/system/auto-add-osd.service
[Unit]
After=network.target

[Service]
ExecStart=/etc/ceph/auto-add-osd.sh

[Install]
WantedBy=default.target
